# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:glassfish4]
command=/bin/bash -c "asadmin start-domain"